# dist/etc/init/garage.conf
description "website to monitor and control the garage door"

start on runlevel [2345]
stop on runlevel [^2345]

chdir /opt/garage
exec python garage.py --port=8888 --db=/var/opt/garage/garage.db --log_file_prefix=/var/log/garage/garage --git=$(cat git_rev)
